# roles/result/files/result.service
[Unit]
Description=Result Web Application
After=network.target

[Service]
WorkingDirectory=/opt/result-app
ExecStart=/usr/bin/node server.js
Environment=POSTGRES_HOST={{ postgresql_host }}
Environment=POSTGRES_PORT=5432
Environment=POSTGRES_USER={{ postgresql_user }}
Environment=POSTGRES_PASSWORD={{ postgresql_password }}
Environment=POSTGRES_DB={{ postgresql_db }}
Restart=always

[Install]
WantedBy=multi-user.target