# roles/result/files/result.service
[Unit]
Description=Result Web Application
After=network.target

[Service]
WorkingDirectory=/opt/result-app
ExecStart=/usr/bin/node server.js
Environment=DB_HOST={{ postgresql_host }}
Environment=DB_PORT=5432
Environment=DB_USER={{ postgresql_user }}
Environment=DB_PASSWORD={{ postgresql_password }}
Environment=DB_NAME={{ postgresql_db }}
Restart=always

[Install]
WantedBy=multi-user.target